# root/opt/config/plexreport_schedule.cron
# Sends a weekly Plex Report email on Friday's at 10:30hrs. See http://www.adminschoice.com/crontab-quick-reference for cron help.
30 10 * * 5 docker exec plexReport plexreport &> /dev/null

