# root/opt/nowshowing_schedule.cron
# NowShowing schedules for email and web. See http://www.adminschoice.com/crontab-quick-reference for cron help.
# Email report schedule for Friday's at 10:30am
30 10 * * 5 emailreport &> /dev/null
# Webpage report schedule every day at 11:30pm
30 23 * * * webreport &> /dev/null
